# tmxaudit/rpmbuild/tmxaudit-server/tmxaudit-server/etc/cron.daily/tmxaudit
#!/bin/bash
# cron.daily/tmxaudit

# Clean the client host backups:
/u01/app/adminscripts/clientbu-reaper.sh &> /tmp/clientbu-reaper.log
