# tmxaudit/rpmbuild/tmxaudit-client/tmxaudit-client/etc/cron.d/killvnc
20 23 * * 5 root /u01/app/adminscripts/killvnc.pl &> /dev/null
