# tmxaudit/rpmbuild/tmxaudit-client/tmxaudit-client/etc/cron.d/timesync
0 */4 * * * root /u01/app/adminscripts/timesync.sh &> /dev/null
