# tmxaudit/rpmbuild/tmxaudit-client/tmxaudit-client/etc/cron.d/tmxaudit
0 10 * * * root /usr/bin/yum -y update tmxaudit-client &> /dev/null
