# tmxaudit/rpmbuild/tmxaudit-client/tmxaudit-client/etc/cron.d/tmxbu-agent
15 1 * * * root /u01/app/adminscripts/tmxbu-agent.sh /u01/app/adminscripts/bu.list  &> /dev/null
