# tmxaudit/rpmbuild/tmxaudit-client/solaris/tmxaudit-client-3.1.2-6-install/etc/cron.d/authclient
0 21 * * * root /u01/app/auth-client/authclient.pl &> /dev/null
